# lib/backend/readHgData.conf
# This is the upstart config to run the shell script as service
console output

start on runlevel [2345]
stop on runlevel [!2345]
#start on startup
respawn
pre-start script

  # Perform whatever checks you like here (maybe checking
  # '/etc/default/foo' to see if the service is enabled # or not).
  #
  # if there are no problems detected, simply "exit 0", else do
  # something like this...

  # display an error message to stderr *on the console* and also write
  # the same message to the system log.
  logger -is -t "$UPSTART_JOB" "ERROR: foo!"

  # tell Upstart not to start the main process for the job.
#  exit 1
end script

# correct this path to your needs
exec /opt/ra5/databases/readHgData.sh